# roles/nash/files/etc/profile.d/nash.sh
# NASH
NASHROOT="/home/lucianofaustino/.local/share/nash"
NASHPATH="/home/lucianofaustino/.config/nash"
export NASHROOT
export NASHPATH
