# ansible/roles/nash/templates/etc/profile.d/nash.sh.j2
# NASH
NASHPATH="{{ nash_path }}"
NASHROOT="{{ nash_root }}"
export NASHROOT
export NASHPATH
